# bootstrap/var/service/ftpd/run
#!/data/data/pl.sviete.dom/files/usr/bin/sh
exec tcpsvd -vE 0.0.0.0 8021 ftpd /data/data/pl.sviete.dom/files/home
